# stage2/06-usb-gadget/files/etc/systemd/system/rpi-cam-gadget-fallback.service
[Unit]
Description=Fall back from UVC to NCM if no USB video stream was opened
# Same OTG gate as the gadget setup; no-op on non-OTG boards.
ConditionPathExists=/run/rpi-cam-gadget.enabled
# Ordering only: we want rpi-camera up first so the pump has had a chance to
# flag an opened UVC stream before we decide.
After=rpi-cam-gadget.service rpi-camera.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/rpi-cam-gadget-mode.sh fallback
