# stage2/06-usb-gadget/files/etc/systemd/system/rpi-cam-gadget-fallback.timer
[Unit]
Description=Grace period before UVC->NCM gadget fallback
# No ConditionPathExists here: the timer is armed by timers.target early in
# boot, possibly before rpi-cam-gadget-detect.service writes the .enabled flag,
# which would skip the timer entirely. The gate lives on the *service* it
# triggers (which is ordered After= detect), so a non-OTG board still no-ops.

[Timer]
# 30 s after boot: long enough for rpi-camera to start and for a host to open
# the UVC stream, short enough that an unused webcam flips to networking fast.
OnBootSec=30
AccuracySec=1s

[Install]
WantedBy=timers.target
