# stage2/06-usb-gadget/files/etc/systemd/system/NetworkManager.service.d/10-dnsmasq-chown.conf
# NetworkManager's shared-mode dnsmasq (used for the usb0 gadget network)
# drops to an unprivileged user and then cannot chown its PID file,
# logging on every connect:
#   dnsmasq: chown of PID file /run/nm-dnsmasq-usb0.pid failed:
#   please add capability CAP_CHOWN
# Add CAP_CHOWN to NetworkManager's ambient set so the dnsmasq child it
# spawns keeps the capability across its setuid. NM itself runs as root
# and is unaffected; this only changes what its non-root children inherit.
# Purely cosmetic — DHCP on usb0 works either way.
[Service]
AmbientCapabilities=CAP_CHOWN
